# static/etc/services.d/plex_requests/run
#!/usr/bin/with-contenv bash
umask 0002

export HOME=/app
cd /app/plex-requests

echo "Starting Plex Requests.. This can take a while.. Make some coffee."
exec s6-setuidgid bytesized meteor
